# src/supervisor/cunxun-ui.ini
[program:cunxun-ui]
command=/root/www/work/gopath/bin/caddy -conf Caddyfile
directory=/data/apps/html/wangqingang/cunxun-ui/release/
autostart=true
autorestart=true
startsecs=2
startretries=3
stdout_logfile=/data/logs/supervisor/cunxun-ui/access.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=20
stderr_logfile=/data/logs/supervisor/cunxun-ui/stderr.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=2
